# systemd/.config/systemd/user/study-focus-block.target
[Unit]
Description=Study Focus Application Blocker
Documentation=man:systemd.special(7)
StopWhenUnneeded=yes
