# systemd/.config/systemd/user/yapyap.service
[Unit]
Description=Yapyap voice dictation (block paste mode)

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/user/yapyap.env
Environment=MODEL=tiny.en-q8_0
Environment=YDOTOOL_SOCKET=%t/ydotool.socket
ExecStart=/usr/bin/bash -lc 'exec /home/joebutler/.local/bin/yapyap | while IFS= read -r l; do /usr/bin/wl-copy -- "$l"; /usr/bin/ydotool key 29:1 47:1 47:0 29:0; done'
Restart=always
RestartSec=1
Nice=10

[Install]
WantedBy=default.target
